# .config/systemd/user/service-watcher@.path
[Path]
PathChanged=%h/.config/%i

[Install]
WantedBy=graphical-session.target
Requires=%i.service
After=%i.service
